# rc.d_halter
#!/bin/sh
 
PATH=/sbin:/usr/sbin:/bin:/usr/bin:/usr/local/bin

. /etc/rc.subr
 
name=halter
rcvar=${name}_enable
 
app_dir="/opt/${name}"
command="/usr/sbin/daemon"
command_args="-p /var/run/${name}.pid ${app_dir}/${name}.py"
 
load_rc_config $name
run_rc_command "$1"
